# sleep.target
[Unit]
Description=User level sleep target
StopWhenUnneeded=yes
RemainAfterExit=yes
